# patches/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
console::askfirst:-${LOGIN}
